# supervisord.conf
[supervisord]
nodaemon=true

[program:squid3]
command=/usr/sbin/squid3 -N
